# var/spool/cron/crontabs/pi
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.d5PmZq/crontab installed on Mon Jan 25 16:00:41 2016)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
# m h  dom mon dow   command
30 5 * * * sudo tvservice -p && sudo pkill lightdm
30 7 * * * sudo tvservice -o && pkill iceweasel
30 18 * * * sudo tvservice -p && sudo pkill lightdm
30 20 * * * sudo tvservice -o && pkill iceweasel
